# client/crontab
# commit monitor data
*/1 * * * * /usr/bin/python /opt/tools/monitor.py